# dataset_evaluation/dataset-evaluation.service
[Unit]
Description=GCP Pub/Sub to dataset evaluation
After=network.target

[Service]
ExecStart=/opt/dataset_evaluation/venv/bin/python3 /opt/dataset_evaluation/main.py
WorkingDirectory=/opt/dataset_evaluation/
Restart=always
RestartSec=10
User=serviceuser

[Install]
WantedBy=multi-user.target
